# scripts/install/systemd/nqrust-manager.service
[Unit]
Description=NQRust-MicroVM Manager
Documentation=https://github.com/NexusQuantum/NQRust-MicroVM
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=nqrust
Group=nqrust
# Add docker group for container feature (pull/manage Docker images)
SupplementaryGroups=docker
WorkingDirectory=/opt/nqrust-microvm

# Load environment from file
EnvironmentFile=/etc/nqrust-microvm/manager.env

# Manager binary
ExecStart=/opt/nqrust-microvm/bin/manager

# Restart policy
Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nqrust-manager

# Manager needs to run sudo for guest agent installation (mount rootfs)
# Therefore we cannot use restrictive security settings like ProtectSystem=strict
# The manager runs as unprivileged 'nqrust' user, security comes from that
NoNewPrivileges=false

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Allow core dumps for debugging (optional)
# LimitCORE=infinity

[Install]
WantedBy=multi-user.target
